# scripts/systemd/backup.timer
[Unit]
Description=Daily timer for plex_playlist backups

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
